# setup/dns.conf
[program:dnsproxy]
user = nobody
command = /usr/bin/python /dns-proxy/server.py
